# templates/usr/lib/systemd/system/carbon-cache.service.erb
# /usr/lib/systemd/system/carbon-cache-<%= @title %>.service
[Unit]
Description=Graphite Carbon Cache <%= @title %>
After=network.target

[Service]
Type=forking
StandardOutput=syslog
StandardError=syslog
ExecStart=/usr/bin/carbon-cache --config=/etc/carbon/carbon.conf --pidfile=/var/run/carbon-cache-<%= @title %>.pid --instance=<%= @title %> start
ExecReload=/bin/kill -USR1 $MAINPID
PIDFile=/var/run/carbon-cache-<%= @title %>.pid

[Install]
WantedBy=multi-user.target