# templates/usr/lib/systemd/system/carbon-cache.service.erb
[Unit]
Description=Graphite Carbon Cache <%= @cache %>
After=network.target

[Service]
Type=forking
StandardOutput=syslog
StandardError=syslog
ExecStart=/usr/bin/carbon-cache --config=<%= @config_dir %>/<%= @config_filename %> --instance=<%= @cache %> --pidfile=<%= @pid_dir %>/carbon-cache-<%= @cache %>.pid --logdir=<%= @log_dir %>/ start
ExecReload=/bin/kill -USR1 $MAINPID
PIDFile=<%= @pid_dir %>/carbon-cache-<%= @cache %>.pid

[Install]
WantedBy=multi-user.target